# vcontrold/rootfs/etc/services.d/vclient_pub/run
#!/usr/bin/with-contenv bashio
source /etc/services.d/get_mqtt_settings.sh
source /etc/services.d/get_vcontrold_settings.sh
sleep 5
while true; do
    bashio::log.info "Executing vclient_pub command..."
    if flock -w 120 /tmp/vcontrol.lock vclient -h "$VCONTROL_HOST" -p "$VCONTROL_PORT" -f /etc/vcontrold/1_mqtt_commands.txt -t /etc/vcontrold/2_mqtt.tmpl -x /etc/vcontrold/3_mqtt_pub.sh; then
        bashio::log.info "vclient_pub executed successfully."
    else
        bashio::log.error "vclient_pub failed."
    fi
    bashio::log.info "Looping vclient_pub. Sleeping for $REFRESH_RATE seconds."
    sleep "$REFRESH_RATE"
done
