# vcontrold/rootfs/etc/services.d/vcontrold/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the vcontrold service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Declare variables
declare auth_arg
declare cmd
declare command_entry
declare config_commands_raw
declare config_deviceid
declare config_file_vcontrold
declare config_file_vcontrold_source
declare config_tty
declare type
declare value_placeholder
declare -a commands_array=()
declare -r services_dir="${SERVICES_DIR:-/etc/services.d}"
declare -r addon_config_dir="${ADDON_CONFIG_DIR:-/config}"
declare -r legacy_config_dir="${LEGACY_CONFIG_DIR:-/homeassistant/vcontrold}"
declare -r vcontrold_config_dir="${VCONTROLD_CONFIG_DIR:-/etc/vcontrold}"
declare -r config_file_vito="${vcontrold_config_dir}/vito.xml"
declare -r default_vcontrold_config="${vcontrold_config_dir}/vcontrold.xml"
declare -r runtime_dir="${RUNTIME_DIR:-/run/vcontrold}"
declare -r runtime_commands_file="${runtime_dir}/1_mqtt_commands.txt"
declare -r runtime_template_file="${runtime_dir}/2_mqtt.tmpl"
declare -r runtime_output_script="${runtime_dir}/3_mqtt_pub.sh"
declare -r runtime_vcontrold_config="${runtime_dir}/vcontrold.xml"
declare -r tail_bin="${TAIL_BIN:-tail}"
declare -r vcontrold_bin="${VCONTROLD_BIN:-/usr/sbin/vcontrold}"

is_local_vcontrold() {
    case "${VCONTROL_HOST}" in
        localhost|127.0.0.1|::1)
            return 0
            ;;
        *)
            return 1
            ;;
    esac
}

# Fetch MQTT settings
source "${services_dir}/command_utils.sh"
source "${services_dir}/get_mqtt_settings.sh"

# Fetch Vcontrold settings
source "${services_dir}/get_vcontrold_settings.sh"

# Get configuration options from user
config_tty=$(bashio::config 'tty')
config_deviceid=$(bashio::config 'device_id')
config_commands_raw=$(bashio::config 'commands')

# Check if custom vito.xml file exists
if [ -f "${addon_config_dir}/vito.xml" ]; then
    bashio::log.info "Using custom vito.xml file."
    config_file_vito="${addon_config_dir}/vito.xml"
elif [ -f "${legacy_config_dir}/vito.xml" ]; then
    bashio::log.info "Using legacy custom vito.xml file from /homeassistant/vcontrold."
    config_file_vito="${legacy_config_dir}/vito.xml"
else
    bashio::log.info "Using default vito.xml file."
fi

# Check if custom vcontrold.xml file exists
if [ -f "${addon_config_dir}/vcontrold.xml" ]; then
    bashio::log.info "Using custom vcontrold.xml file."
    config_file_vcontrold_source="${addon_config_dir}/vcontrold.xml"
elif [ -f "${legacy_config_dir}/vcontrold.xml" ]; then
    bashio::log.info "Using legacy custom vcontrold.xml file from /homeassistant/vcontrold."
    config_file_vcontrold_source="${legacy_config_dir}/vcontrold.xml"
else
    bashio::log.info "Using default vcontrold.xml file."
    config_file_vcontrold_source="${default_vcontrold_config}"
fi

config_file_vcontrold="${runtime_vcontrold_config}"

# Prepare runtime files
mkdir -p "${runtime_dir}"
rm -f "${runtime_commands_file}" "${runtime_template_file}" "${runtime_output_script}"
printf '#!/bin/sh\n' > "${runtime_template_file}"

# Configure scripts based on commands
load_vcontrold_commands_array "${config_commands_raw}" commands_array
i=0
for command_entry in "${commands_array[@]}"
do
    [[ -n "${command_entry}" ]] || continue

    if [[ "${command_entry}" != *:* ]]; then
        bashio::log.warning "Skipping invalid command entry: ${command_entry}"
        continue
    fi

    ((i=i+1))
    cmd=${command_entry%%:*}
    type=${command_entry#*:}

    if [[ -z "${cmd}" ]]; then
        bashio::log.warning "Skipping empty command entry at position ${i}."
        ((i=i-1))
        continue
    fi

    bashio::log.info "Configuring command ${i}: ${cmd} with type ${type}"
    printf '%s\n' "${cmd}" >> "${runtime_commands_file}"

    value_placeholder="\$${i}"
    if [[ "${type}" == "STRING" ]]; then
        value_placeholder="\$R${i}"
    fi

    printf 'if [ "x$E%d" = "xOK" ]; then\n' "${i}" >> "${runtime_template_file}"
    printf '  mosquitto_pub -h %q -p %q' "${MQTT_HOST}" "${MQTT_PORT}" >> "${runtime_template_file}"
    for auth_arg in "${MQTT_AUTH_ARGS[@]}"; do
        printf ' %q' "${auth_arg}" >> "${runtime_template_file}"
    done
    printf ' -t %q/$C%d -m "%s"\n' "${MQTT_TOPIC}" "${i}" "${value_placeholder}" >> "${runtime_template_file}"
    printf 'fi\n' >> "${runtime_template_file}"
done

if [[ ! -s "${runtime_commands_file}" ]]; then
    bashio::exit.nok "No valid commands configured."
fi

cp "${config_file_vcontrold_source}" "${config_file_vcontrold}"

sed -i "s/#DEVICEID#/${config_deviceid}/g" "${config_file_vcontrold}"
sed -i "s|#VITOXML#|${config_file_vito}|g" "${config_file_vcontrold}"

if bashio::config.true 'debug'; then
    sed -i "s/#DEBUG#/y/g" "${config_file_vcontrold}"
else
    sed -i "s/#DEBUG#/n/g" "${config_file_vcontrold}"
fi

# Keep custom vcontrold.xml files compatible when they still reference legacy or relative vito.xml paths.
sed -i "s|href=\"/config/vcontrold/vito.xml\"|href=\"${config_file_vito}\"|g" "${config_file_vcontrold}"
sed -i "s|href=\"/homeassistant/vcontrold/vito.xml\"|href=\"${config_file_vito}\"|g" "${config_file_vcontrold}"
sed -i "s|href=\"./vito.xml\"|href=\"${config_file_vito}\"|g" "${config_file_vcontrold}"
sed -i "s|href=\"vcontrold/vito.xml\"|href=\"${config_file_vito}\"|g" "${config_file_vcontrold}"
sed -i "s|href=\"vito.xml\"|href=\"${config_file_vito}\"|g" "${config_file_vcontrold}"

# Start the vcontrold service or a dummy process for external control
if is_local_vcontrold; then
    bashio::log.info "Starting local vcontrold..."

    if [[ -z "${config_tty}" ]]; then
        bashio::exit.nok "A tty device is required when running local vcontrold."
    fi

    exec "${vcontrold_bin}" -x "${config_file_vcontrold}" -n -U root -d "${config_tty}"
fi

bashio::log.info "Using remote vcontrold at ${VCONTROL_HOST}:${VCONTROL_PORT}. Keeping process alive."
exec "${tail_bin}" -f /dev/null
