# vcontrold/rootfs/etc/services.d/vclient_sub/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the MQTT-based vclient_sub service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

source /etc/services.d/get_mqtt_settings.sh
source /etc/services.d/get_vcontrold_settings.sh

declare command
declare payload
declare request
declare topic
declare value
declare -r lock_file="/tmp/vcontrol.lock"

if [[ -z "$MQTT_HOST" || -z "$MQTT_PORT" || -z "$MQTT_TOPIC" || -z "$REFRESH_RATE" || -z "$VCONTROL_HOST" || -z "$VCONTROL_PORT" ]]; then
    bashio::log.error "Missing required MQTT or Vcontrold environment variables."
    exit 1
fi

bashio::log.info "Starting MQTT subscription to topic: ${MQTT_TOPIC}/#"

while true; do
    mosquitto_sub -v -h "$MQTT_HOST" -p "$MQTT_PORT" "${MQTT_AUTH_ARGS[@]}" -t "${MQTT_TOPIC}/#" | while IFS= read -r payload; do
        topic=${payload%% *}

        if [[ "${payload}" == *" "* ]]; then
            value=${payload#* }
        else
            value=""
        fi

        command=${topic#"$MQTT_TOPIC"/}
        if [[ "${command}" == "${topic}" || "${command}" != set* ]]; then
            continue
        fi

        bashio::log.info "Received command: [${command}] with value: ${value}"

        if [[ -n "${value}" ]]; then
            request="${command} ${value}"
        else
            request="${command}"
        fi

        if flock -w 120 "${lock_file}" vclient -h "$VCONTROL_HOST" -p "$VCONTROL_PORT" -o /dev/stdout -c "${request}"; then
            bashio::log.info "Command [${command}] executed successfully."
        else
            bashio::log.error "Failed to execute command [${command}] with value: ${value}"
        fi
    done

    bashio::log.warning "MQTT subscription ended. Reconnecting in ${REFRESH_RATE} seconds."
    sleep "$REFRESH_RATE"
done
