# vcontrold/rootfs/etc/services.d/vclient_sub/run
#!/usr/bin/with-contenv bashio
source /etc/services.d/get_mqtt_settings.sh
source /etc/services.d/get_vcontrold_settings.sh
bashio::log.info "Starting MQTT subscription to topic: $MQTT_TOPIC/#"
while true; do
    mosquitto_sub -v -h "$MQTT_HOST" -p "$MQTT_PORT" -u "$MQTT_USER" -P "$MQTT_PASSWORD" -t "$MQTT_TOPIC/#" | while read -r payload; do
        topic=$(echo "$payload" | cut -d' ' -f1)
        value=$(echo "$payload" | cut -d' ' -f2)
        if [[ "$topic" =~ "/set" ]]; then
            command=$(echo "$topic" | cut -d'/' -f2)
            get_command=$(echo "$command" | sed 's/^set/get/')
            bashio::log.info "Sending command: [${command}]: ${value}"
            flock -w 120 /tmp/vcontrol.lock bash -c "
                vclient -h '$VCONTROL_HOST' -p '$VCONTROL_PORT' -o /dev/stdout -c '${command} ${value}' && \
                sleep 1 && \
                result=\$(vclient -h '$VCONTROL_HOST' -p '$VCONTROL_PORT' -o /dev/stdout -c '${get_command}') && \
                echo \"\$result\" | mosquitto_pub -h '$MQTT_HOST' -p '$MQTT_PORT' -u '$MQTT_USER' -P '$MQTT_PASSWORD' -t '$MQTT_TOPIC/${get_command}' -l
            "
            bashio::log.info "Command [${command}] executed, result published to MQTT."
        fi
    done
    sleep 1
done
